# deploy/systemd/paper-watch.service
[Unit]
Description=paper-watch digest run

[Service]
Type=oneshot
# %h expands to the user's home, so these units work for any user without
# hardcoding a username.
WorkingDirectory=%h/Code/paper-watch
# uv lives in the user's ~/.local/bin, which is not on the PATH that systemd
# uses, so call it by absolute path.
ExecStart=%h/.local/bin/uv run paper-watch run
# Output goes to the journal: journalctl --user -u paper-watch
